# docker/supervisor/conf/socket.conf
[program:socket]
process_name=%(program_name)s
command=php /usr/src/app/tcp-socket.php
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
